# mqtt-shell-agent.service
[Unit]
Description=MQTT Remote Shell Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mqttshell
Group=mqttshell
WorkingDirectory=/opt/mqtt-shell

# Environment variables for credentials
Environment="MQTT_HOST=broker.example.com"
Environment="MQTT_PORT=1883"
Environment="MQTT_USER=yourusername"
Environment="MQTT_PASSWORD=your-broker-password"
Environment="MQTT_AGENT_PASSWORD=unique-agent-secret"
Environment="MQTT_ID=agent-01"

ExecStart=/usr/bin/python3 /opt/mqtt-shell/mqtt_shell_agent.py

Restart=always
RestartSec=3

# Hardening (optional; relax these if you need full login behavior)
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=off
ProtectHome=false
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
